# hosts/pajlada-combine-hosts.service
[Unit]
Description=Combine hosts

[Service]
ExecStart=/opt/pajlada-dotfiles/combine-hosts.sh
